# packages/linux-firmware/brcmfmac_sdio-firmware/system.d/brcmfmac-firmware.service
[Unit]
Description=Broadcom SDIO Firmware Symlink Service
After=kernel-overlays.service
Before=ludo.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/usr/bin/brcmfmac-firmware-setup"

[Install]
WantedBy=local-fs.target
