# packages/x11/xserver/xorg-server/system.d/xorg-configure-modesetting.service
[Unit]
Description=configure Xorg Server for modesetting

ConditionKernelCommandLine=!installer

[Service]
Type=oneshot
ExecStart=/usr/lib/xorg/xorg-configure modesetting
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
